# templates/etc/systemd/system/cockpit.socket.d/listen.conf.erb
[Socket]
# This is not a typo, it's how systemd resets the ListenStream setting
ListenStream=
ListenStream=<%= scope['cockpit::port'] %>
